# .config/systemd/user/bashbullet.service
[Unit]
Description=BashBullet

[Service]
ExecStart=/usr/sbin/bashbullet

[Install]
WantedBy=default.target
