# .config/systemd/user/gpaste.service
[Unit]
Description=Gpaste

[Service]
ExecStart=/usr/bin/gpaste-client daemon

[Install]
WantedBy=default.target
