# .config/systemd/user/startx.service
[Unit]
Description='startx'
After=getty@tty1.service

[Service]
Type=simple
ExecStart=/usr/bin/startx

[Install]
WantedBy=getty@tty1.service
